# playbook/files/bucky.init.conf
description "Bucky stats logger"

start on startup

env USER=nobody
env GROUP=nobody

script
  /usr/bin/bucky /etc/bucky/bucky.conf
end script
